# config/supervisor/oauth2_proxy.conf
[program:oauth2_proxy]
command=/usr/local/bin/oauth2_proxy %(ENV_OAUTH2_PROXY_ARGS)s -http-address=0.0.0.0:80 -upstream=file:///home/shared/s3/#/
