# config/crontab
* * * * * /usr/bin/python /home/scripts/run_s3sync.py
